# packaging/rootfs/overlay/etc/systemd/system/kea-dhcp4-server.service.d/override.conf
[Service]
# Avoid lock directory under /run/lock in unprivileged LXC where run-lock.mount can fail.
Environment="KEA_LOCKFILE_DIR=/run/kea"
RuntimeDirectory=
RuntimeDirectory=kea
RuntimeDirectoryMode=0750
